# sh/profile.d/firefox.sh
pathappend /opt/firefox/                PATH

export PATH
